# conf/supervisor.conf
[program:apache2]
command=/usr/local/bin/apache2-foreground
redirect_stderr=true

[program:worker]
command=/var/www/MISP/app/Console/worker/start.sh
user=www-data
startsecs=0
redirect_stderr=true